# backend/airqmon-alerter.service
[Unit]
Description=AirQMon Alerter Service
After=network.target
PartOf=airqmon.target

[Service]
Type=simple
User=admin
EnvironmentFile=-/etc/airqmon/airqmon.env
ExecStart=/usr/local/bin/airqmon-launch alerter
Restart=on-failure

[Install]
WantedBy=multi-user.target
